# c6u/files/etc/init.d/z_setup_first_stage
#!/bin/sh /etc/rc.common

START=99
LOGGER_PROMPT="LiftOff"
. /root/updater/logger.sh

start() {
    log "Lift off complete.."
    log "Preparing first stage.."

    chmod +x /root/setup/first_stage.sh
    /root/setup/first_stage.sh
}
